# com.euro-quote.widget.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.euro-quote.widget</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/yourusername/widget-macos-wise-eur-to-brl/widget.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/yourusername/widget-macos-wise-eur-to-brl</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/euro-quote.widget.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/euro-quote.widget.out</string>
</dict>
</plist> 